# admin-scripts/systemd/server-watchdog.service
[Unit]
Description=Server Watchdog (self-healing checks)
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/kv1/telegram.env
ExecStart=/home/homepage/admin-scripts/server-watchdog.sh
